# roles/encrypted_volumes/templates/protected-services.target.j2
[Unit]
Requires=mount-encrypted-binds.target
{% for u in allunits -%}
Requires={{ u }}
{% endfor %}
